# subworkers/plists/com.elia.tempack-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.elia.tempack-dev</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/vakandi/EliaAI/subworkers/scripts/trigger_tempack_dev.sh</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>Disabled</key>
    <true/>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/vakandi/.opencode/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/vakandi</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/vakandi/EliaAI</string>

    <key>StandardOutPath</key>
    <string>/Users/vakandi/EliaAI/subworkers/logs/tempack_dev.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vakandi/EliaAI/subworkers/logs/tempack_dev.log</string>
</dict>
</plist>
